# src/ansible/kube_eng/roles/cloud_provider_mdns_configuration/templates/org.mrmat.kubeeng.cloudprovidermdns.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.mrmat.kubeeng.cloudprovidermdns</string>
	<key>ProgramArguments</key>
	<array>
        <string>{{ tool_cloud_provider_mdns }}</string>
	</array>
    <key>WorkingDirectory</key>
    <string>{{ directory }}</string>
    <key>StandardErrorPath</key>
    <string>{{ directory }}/cloud-provider-mdns.log</string>
    <key>StandardOutPath</key>
    <string>{{ directory }}/cloud-provider-mdns.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>UserName</key>
    <string>{{ user }}</string>
</dict>
</plist>
